# bamako-sqld.service
[Unit]
Description=Bamako sqld (libSQL server)
After=network.target

[Service]
Type=simple
User=bamako
Group=bamako
WorkingDirectory=/opt/bamako

ExecStart=/opt/bamako/sqld \
    --http-listen-addr=0.0.0.0:8092 \
    --db-path=/opt/bamako/data

EnvironmentFile=/opt/bamako/.env

Restart=on-failure
RestartSec=10s
StartLimitInterval=5min
StartLimitBurst=3

LimitNOFILE=65536
LimitNPROC=4096

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/bamako/data
ReadWritePaths=/opt/bamako/logs

StandardOutput=journal
StandardError=journal
SyslogIdentifier=bamako-sqld

[Install]
WantedBy=multi-user.target
